# init/corosync-notifyd.service.in
[Unit]
Description=Corosync Dbus and snmp notifier
Documentation=man:corosync-notifyd
Wants=corosync.service
After=corosync.service

[Service]
EnvironmentFile=-@CONFIGDIR@/corosync-notifyd
ExecStart=@sbindir_expanded@/corosync-notifyd -f $OPTIONS
Type=simple
Restart=on-failure

[Install]
WantedBy=multi-user.target
